# dockerfiles/centos-php-fpm/root/etc/s6/php/run
#!/usr/bin/env bash

exec s6-setuidgid apache php-fpm -F
